# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:mqtt-proxy]
command=/go/src/github.com/ninjablocks/mqtt-proxy/mqtt-proxy -config=/etc/mqtt-proxy/config.toml
redirect_stderr=true