# etc/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:mqtt-proxy]
command=/go/src/github.com/ninjablocks/mqtt-proxy/mqtt-proxy
redirect_stderr=true